# applications/app-meta-transmission/root/etc/uci-defaults/50_app-meta-transmission
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	set transmission.@transmission[-1].web_home=/usr/share/transmission/web
	commit transmission
EOF

exit 0