# applications/app-meta-transmission/root/etc/uci-defaults/50_app-meta-transmission
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete transmission.@transmission[-1].web_home
	commit transmission
EOF

sed -i 's/reload_service(/reload_service0(/g' /etc/init.d/transmission

exit 0
