# applications/app-meta-aria2/root/etc/uci-defaults/50_app-meta-aria2
#!/bin/sh

# fix ca_certificate
uci -q get aria2.main.ca_certificate >/dev/null && exit 0

uci -q batch <<-EOF >/dev/null
	set aria2.main.ca_certificate=/etc/ssl/certs/ca-certificates.crt
	commit aria2
EOF

exit 0
